# apps/auto-mcs.conf
[program:auto-mcs]
command=/bad/auto-mcs
autorestart=true
