# apps/fluxbox.conf
[program:fluxbox]
environment=DISPLAY=:0
command=fluxbox
autorestart=true