# apps/ratpoison.conf
[program:ratpoison]
command=/usr/bin/ratpoison --display :0
autorestart=true