# apps/software.conf
[program:software]
command=/bad/auto-mcs
autorestart=true