# apps/software.conf
[program:software]
environment=DISPLAY=:0
command=/bad/auto-mcs
autorestart=true