# apps/websockify.conf
[program:websockify]
command=websockify --key /bad/p.key --cert /bad/c.crt --ssl-only --web /bad/novnc 80 localhost:5900
autorestart=true