# apps/websockify.conf
[program:websockify]
command=websockify --web /bad/novnc 80 localhost:5900
autorestart=true