# apps/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 900x900x24 -listen tcp -ac
autorestart=true
